# gitlab/files/default/gitlab
#!/bin/sh

# chkconfig:   345 99 1
# description: gitlab

/usr/bin/gitlab-ctl $@
